# scripts/rmd.service
[Unit]
Description=rmd

[Service]
User=root
Group=root
ExecStart=/usr/bin/rmd

[Install]
WantedBy=multi-user.target
